# etc/supervisord.d/saslauthd-krb5.ini
[program:saslauthd]
command=/usr/sbin/saslauthd -a kerberos5 -c -m /var/spool/postfix/var/run/saslauthd -n 5
pidfile=/var/spool/postfix/var/run/saslauthd/saslauthd.pid
