# etc/supervisord.d/saslauthd-pam.ini
[program:saslauthd]
command=/usr/sbin/saslauthd -a pam -c -m /var/spool/postfix/var/run/saslauthd -n 5
pidfile=/var/spool/postfix/var/run/saslauthd/saslauthd.pid
